# jboss-eap/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:jboss]
command=$JBOSS_HOME/bin/standalone.sh --debug -bmanagement 0.0.0.0 -b 0.0.0.0

